# rails_root/script/crontab
MAILTO=gigix1980@gmail.com
45 23 * * * bash -lc "cd /home/ohfy/ohfy/rails_root && /home/ohfy/.gems/bin/bundle exec rails runner -e production script/share.rb"