# rails_root/script/crontab
* * * * * cd /home/ohfy/ohfy/rails_root && bundle exec rails runner -e production script/share.rb